# project/deploy/lcd1602.service
[Unit]
Description=LCD1602 Display Service

[Service]
ExecStart=/usr/bin/lcd1602_app
Restart=always

[Install]
WantedBy=multi-user.target
